# conf.d/caffe_jupyter.conf
[program:caffe]
command=/root/startjupyter.sh 8888
directory=${HOME}
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
killasgroup=true
