# conf.d/spark.conf
[program:spark_ipython]
command=/root/startspark.sh
directory=/opt/spark-2.1.0-bin-hadoop2.7/bin/
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true
